# jupyter.conf
[program:jupyter]
command=[PATH]
directory=[DIR]
user=ubuntu
numprocs=1
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 600
stopasgroup=true
stopsignal=QUIT
stderr_logfile=[LOG_FILE_DIR]err.log
stdout_logfile=[LOG_FILE_DIR]out.log
