# jupyter.conf
[program:jupyter]
command=[PATH_TO_EXEC]
directory=[NOTEBOOK_DIR]
user=ubuntu
numprocs=1
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 600
stopasgroup=true
stopsignal=QUIT
stderr_logfile=[LOG_DIR]/jupyter.err.log
stdout_logfile=[LOG_DIR]/jupyter.out.log
